# eos-phone-home.path
[Unit]
Description=Send system activation/daily ping messages to Endless
Wants=network.target

[Path]
DirectoryNotEmpty=/home
Unit=eos-phone-home.timer

[Install]
WantedBy=multi-user.target
